# dot_config/systemd/user/private_restic-backup.timer
[Unit]
Description=Schedule for Restic Backup

[Timer]
OnCalendar=hourly
Unit=restic-backup.service

[Install]
WantedBy=timers.target
